# crontab
#soal1
14 14 * 2 5 ~/SoalShift_modul1_B02/soal1.sh

#soal4
0 */1 * * * /bin/bash ~/SoalShift_modul1_B02/soal4.sh

#soal5
2-30/6 * * * * ~/SoalShift_modul1_B02/soal5.sh
